# releng/airootfs/etc/profile.d/firstboot.sh
#!/bin/bash
username=$(whoami)
/home/$username/dotfiles/install.sh
